# dist/test/edge-agent/edge-agent.service
[Unit]
Description=Edge IDPS Agent
After=network.target

[Service]
EnvironmentFile=/etc/default/edge-agent
ExecStart=/usr/local/bin/edge-agent
Restart=always
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
